# rc.local
#!/usr/local/bin/bash
echo
echo "========== startup forever =============="
cd /usr/local/etc
/usr/local/bin/forever start server.js
